# etc/supervisor/conf.d/php-fpm.conf
[program:php7.2]
command=/usr/sbin/php-fpm7.2 -F
redirect_stderr=true
numprocs=1
autostart=true
autorestart=true
priority=10
